# bin/apps/aspiesoft-clamav-scanner/aspiesoft-clamav-download-scanner.service
[Unit]
Description=AspieSoft ClamAV Download Scanner

[Service]
ExecStart=/bin/bash /etc/aspiesoft-clamav-scanner/start.sh

[Install]
WantedBy=multi-user.target
