# bin/apps/empoleos/empoleos.service
[Unit]
Description=Empoleos auto Updates and Backups

[Service]
ExecStart=/bin/bash /etc/empoleos/init.sh

[Install]
WantedBy=multi-user.target
