# bin/apps/empoleos/empoleos.service
[Unit]
Description=AspieSoft Fedora Setup Updates

[Service]
ExecStart=/bin/bash /etc/empoleos/init.sh

[Install]
WantedBy=multi-user.target
